# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stderr
logfile_maxbytes=0
loglevel=info

[program:uvicorn]
command=python -m uvicorn backend.main:app --host 127.0.0.1 --port 8000 --log-level info
directory=/app
autostart=true
autorestart=true
priority=10
startsecs=3
environment=
    PYTHONUNBUFFERED="1",
    GOOGLE_CLOUD_PROJECT="%(ENV_GOOGLE_CLOUD_PROJECT)s",
    FIRESTORE_COLLECTION="%(ENV_FIRESTORE_COLLECTION)s",
    GEMINI_API_KEY="%(ENV_GEMINI_API_KEY)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=bash -c "while ! curl -s http://127.0.0.1:8000/api/health > /dev/null; do sleep 0.5; done && exec nginx -g 'daemon off;'"
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface
